# micro-os/combustion/simple-backup.timer
[Unit]
Description=Daily backup of the system

[Timer]
OnCalendar=daily
RandomizedDelaySec=1h
Persistent=true

[Install]
WantedBy=timers.target‚